# deploy/research_platform.conf
; Supervisor program entry for the research_platform server.
; Install on Emergent deploy box with:
;   sudo cp /app/deploy/research_platform.conf /etc/supervisor/conf.d/research_platform.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start research_platform
; Check:
;   sudo supervisorctl status research_platform
;   tail -f /var/log/supervisor/research_platform.out.log

[program:research_platform]
command=/usr/bin/python3 -u /app/notebooks/research_platform_v6_with_PDF_report.py
directory=/app/notebooks
autostart=true
autorestart=true
startretries=999
startsecs=5
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/var/log/supervisor/research_platform.out.log
stderr_logfile=/var/log/supervisor/research_platform.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
